# binance-collector.service
[Unit]
Description=Binance Orderbook Data Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Vultr VPS部署配置：使用root用户，应用安装在/opt/binance_collector
User=root
Group=root
WorkingDirectory=/opt/binance_collector
ExecStart=/usr/bin/python3 /opt/binance_collector/main.py
ExecStop=/bin/kill -TERM $MAINPID

# 自动重启配置
Restart=always
RestartSec=10

# 资源限制
LimitNOFILE=65536
MemoryLimit=2G

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=binance-collector

# 安全配置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
